# supervisord.conf
[supervisord]
nodaemon=true

[program:seneca]
directory=/usr/src/app
command=npm run-script serve
environment=NODE_ENV=production
autostart=true
autorestart=true
